# deploy/isahn.timer
[Unit]
Description=isahn timer script

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
